# deploy/systemd/mailadmin.service
[Unit]
Description=MailOps Go API
After=network-online.target mysqld.service valkey-mail.service
Wants=network-online.target valkey-mail.service

[Service]
Type=simple
User=mailops
Group=mailops
SupplementaryGroups=valkey-mail
WorkingDirectory=/opt/apps/mailops/backend
EnvironmentFile=/opt/apps/mailops/deploy/systemd/mailadmin.env
ExecStart=/opt/apps/mailops/backend/bin/mailadmin-api
Restart=always
RestartSec=3
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/tmp
UMask=0077

[Install]
WantedBy=multi-user.target
