# deploy/systemd/mailadmin.service.d/redis-unix.conf
[Unit]
After=valkey-mail.service
Wants=valkey-mail.service

[Service]
SupplementaryGroups=valkey-mail dovecot
